# supervisor.conf

[program:snippetbox]
command= ./web
directory=/home/stt/kimo-client/
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
stopsignal=QUIT
numprocs=1
numprocs_start=1
